# systemd/ha-watchdog.service
[Unit]
Description=HA watchdog probe and restart
After=network.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ha-watchdog.sh
