# systemd/mqtt-watchdog.service
[Unit]
Description=MQTT broker watchdog
After=mosquitto.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mqtt-watchdog.sh
