# systemd/victron-tools.service
[Unit]
Description=(Legacy) Victron tools compose — superseded by Dockge + docker-compose.tools.yml
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/user1/victron-ble2mqtt-integration
ExecStart=/usr/bin/docker compose -f docker-compose.tools.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.tools.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
