# systemd/vscode-server-cleanup.service
[Unit]
Description=Clean up stale VS Code Server processes
Wants=vscode-server-cleanup.timer

[Service]
Type=oneshot
User=user1
Group=user1
ExecStart=/usr/local/bin/vscode-server-cleanup.sh
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
